# winstart.service
[Unit]
Description=winStart
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/libinput-gestures-setup start
Restart=always

[Install]
WantedBy=xsession.target
