# winstart.service
[Unit]
Description=winStart
After=xorg.target

[Service]
User=rishabh
ExecStart=/usr/bin/libinput-gestures-setup start
RemainAfterExit=false
StandardOutput=journal
